# deploy/systemd/pi_controller.service
[Unit]
Description=Pi Controller MQTT
After=network.target mosquitto.service

[Service]
Type=simple
User=amine
WorkingDirectory=/home/amine
EnvironmentFile=-/home/amine/.env
ExecStart=/usr/bin/python3 /home/amine/pi_controller.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
